# Linux/SystemdService/httpService.service
[Service]
ExecStart=/home/raji/practice/linux/httpserver.sh
